# etc/jabali-security-web.service
[Unit]
Description=Jabali Security Web Dashboard
After=network.target jabali-security.service
Wants=jabali-security.service

[Service]
Type=simple
ExecStart=/usr/local/jabali-security/venv/bin/python -m daemon web --bind 0.0.0.0 --port 8443
WorkingDirectory=/usr/local/jabali-security
Environment=PYTHONPATH=/usr/local/jabali-security
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jabali-security-web

[Install]
WantedBy=multi-user.target
