# etc/jabali-security.service
[Unit]
Description=Jabali Security Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/jabali-security/venv/bin/python -m daemon start
WorkingDirectory=/usr/local/jabali-security
Environment=PYTHONPATH=/usr/local/jabali-security
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jabali-security

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/security/quarantine /var/lib/jabali-security /var/log/jabali-security /etc/php /etc/ufw /lib/ufw /run
# UFW needs /run/ufw.lock and /proc/sys/net for sysctl tuning
ProtectKernelTunables=no
RuntimeDirectory=jabali-security
ProtectHome=no
PrivateTmp=yes
# CAP_DAC_READ_SEARCH: read any file for scanning
# CAP_DAC_OVERRIDE: delete/move files for quarantine
# CAP_NET_ADMIN: nftables IP blocking
# CAP_KILL: process killer (proactive defense)
# CAP_FOWNER: change permissions on quarantined files
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE CAP_NET_ADMIN CAP_KILL CAP_FOWNER
AmbientCapabilities=CAP_DAC_READ_SEARCH CAP_DAC_OVERRIDE CAP_NET_ADMIN CAP_KILL CAP_FOWNER
LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target
